# deploy/ubuntu/systemd/diaghub-frontend.service
[Unit]
Description=DiagHub Next.js frontend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=diaghub
Group=diaghub
WorkingDirectory=/opt/diaghub/repo/frontend
Environment=NODE_ENV=production
EnvironmentFile=-/etc/diaghub/frontend.env
ExecStart=/usr/bin/env npm run start -- --hostname 127.0.0.1 --port 3000
Restart=on-failure
RestartSec=3
TimeoutStartSec=90
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
